# scripts/systemd/nginx-heartbeat.service
[Unit]
Description=NGINX Heartbeat Report
[Service]
Type=oneshot
ExecStart=/usr/local/bin/nginx-heartbeat.sh
