# scripts/systemd/sysmon.service
[Unit]
Description=System Monitor Collector
After=network.target
[Service]
ExecStart=/usr/local/bin/sysmon-collect.sh
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
